# backend/deploy/systemd/pe-be-close-stale-open-workouts.service
[Unit]
Description=PE Backend close stale open workouts
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/pe-be
EnvironmentFile=-/srv/pe-be/backend/.env.production
TimeoutStartSec=300
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml run --rm -T backend python -m src.jobs.close_stale_open_workouts
