# backend/deploy/systemd/pe-be-postgres-backup.timer
[Unit]
Description=Run PE Backend Postgres backup daily at 02:30

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
