# backend/deploy/systemd/pe-be-close-stale-open-workouts.timer
[Unit]
Description=Run PE Backend stale open workout auto-close daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
